# deploy/drust.service
[Unit]
Description=Drust multi-tenant SQLite BaaS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=drust
Group=drust
WorkingDirectory=/home/kaelsohappy1/tool/drust
EnvironmentFile=/home/kaelsohappy1/tool/drust/.env
ExecStart=/home/kaelsohappy1/tool/drust/target/release/drust
Restart=on-failure
RestartSec=5s
UMask=0077

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/drust /var/log/drust
PrivateTmp=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
# NOTE: MemoryDenyWriteExecute is intentionally NOT set. Edge functions
# (v1.36+) JIT-compile guest wasm via wasmtime's Cranelift backend, which
# must make anonymous memory executable; W^X enforcement makes every
# function upload/invoke fail with "unable to make memory executable
# (EPERM)". The guest sandbox is enforced inside wasmtime (epoch deadline +
# ResourceLimiter memory cap + locked-down WasiCtx), not by this directive.

[Install]
WantedBy=multi-user.target
