# deploy/drust-janitor.timer
[Unit]
Description=Daily Drust trash janitor at 04:00

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
